# live-build/config/includes.chroot/etc/systemd/system/gdm.service.d/wait-for-live-config.conf
[Unit]
After=live-config.service
Wants=live-config.service
